# deploy/aneo-led.service
[Unit]
# Start after SPI device is available
After=dev-spidev0.0.device
Wants=dev-spidev0.0.device

# Start very early in boot
Before=multi-user.target
DefaultDependencies=no

[Service]
# Keep the service simple and restart on failure
Type=simple
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aneo-led
ExecStart=/usr/bin/env python3 /opt/aneo_led/led_aneo.py
WorkingDirectory=/opt/aneo_led
User=root
Group=root

[Install]
WantedBy=sysinit.target
