# deploy/jetson-led.service
[Unit]
Description=Jetson Orin Nano LED CPU Monitor
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
